# recipes-multimedia/mtk-vpud/files/vpud.service
[Unit]
Description=MediaTek VPUD daemon

[Service]
ExecStart=/usr/bin/vpud
Restart=always

[Install]
WantedBy=multi-user.target
